# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ollama]
command=ollama serve
environment=OLLAMA_HOST="0.0.0.0:11434",OLLAMA_MODELS="/workspace/models/ollama"
autostart=true
autorestart=true
stdout_logfile=/var/log/ollama.log
stderr_logfile=/var/log/ollama_err.log
priority=10

[program:comfyui]
command=python3 -u /app/ComfyUI/main.py --port 8199 --listen 0.0.0.0 --reserve-vram 4 --disable-cuda-malloc --enable-cors-header * --preview-method none --disable-auto-launch
directory=/app/ComfyUI
environment=PYTHONPATH="/app/ComfyUI",PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/comfyui.log
stderr_logfile=/var/log/comfyui_err.log
priority=20

[program:backend]
command=python3 -u /app/backend/server.py
directory=/app/backend
environment=PYTHONPATH="/app/backend:/app/ComfyUI",PYTHONUNBUFFERED="1",CORS_ORIGINS="*",COMFY_VIEW_BASE="/comfy"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend_err.log
priority=30

[program:nginx]
command=nginx
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx_access.log
stderr_logfile=/var/log/nginx_error.log
priority=40

